# gh-notifier.service
[Unit]
Description=GitHub Notifier Service
After=network.target

[Service]
Type=simple
ExecStart=/path/to/your/compiled/binary
Environment="GITHUB_TOKEN=your_github_token"
Restart=on-failure

[Install]
WantedBy=default.target
